# deploy/supervisor.conf
[program:mcptools]
command=/var/www/mcptools/.venv/bin/gunicorn main:app -c /var/www/mcptools/deploy/gunicorn.conf.py
directory=/var/www/mcptools
user=www-data
group=www-data
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/mcptools/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/mcptools/supervisor_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

# Environment variables
environment=PYTHONPATH="/var/www/mcptools"

# Process priority
priority=999

# Stop signals
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
